# telegramBot.service
[Unit]
Description=Bazi Telegram Bot Agent Service (Rust)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=dietpi
Group=dietpi
WorkingDirectory=/home/dietpi/BaziAgentWorkflow
ExecStart=/home/dietpi/BaziAgentWorkflow/target/release/BaziAgentWorkflow
Restart=always
RestartSec=5
# Environment file
EnvironmentFile=/home/dietpi/BaziAgentWorkflow/.env

[Install]
WantedBy=multi-user.target