# DifyTelegramBot/telegramBot.service
[Unit]
Description=Dify Bazi Bot Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=dietpi
Group=dietpi
WorkingDirectory=/home/dietpi/DifyBaziTelegramBot
ExecStart=/home/dietpi/DifyBaziTelegramBot/.venv/bin/python3 /home/dietpi/DifyBaziTelegramBot/telegramBot.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target